# goes_timelapse/rootfs/etc/services.d/goes-timelapse/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

export GOES_DATA_DIR="$(bashio::config 'data_dir')"
state_dir="$(bashio::config 'state_dir')"
if [[ -n "$state_dir" && "$state_dir" != "null" ]]; then
  export GOES_STATE_DIR="$state_dir"
fi
scratch_dir="$(bashio::config 'scratch_dir')"
if [[ -n "$scratch_dir" && "$scratch_dir" != "null" ]]; then
  export GOES_SCRATCH_DIR="$scratch_dir"
fi
export GOES_POLL_MINUTES="$(bashio::config 'poll_minutes')"
export GOES_FRAME_COUNT="$(bashio::config 'frame_count')"
export GOES_GIF_FPS="$(bashio::config 'gif_fps')"
export GOES_RAW_HISTORY="$(bashio::config 'raw_history')"
export GOES_SOLAR_MARGIN_HOURS="$(bashio::config 'solar_margin_hours')"
export GOES_LOG_LEVEL="$(bashio::config 'log_level')"

exec python3 -m goes_timelapse
